# pypiserver/files/upstart.conf
{% from "pypiserver/map.jinja" import pypiserver with context %}

description "{{ pypiserver.description }}"
start on runlevel [2345]
stop on runlevel [06]

exec {{ pypiserver.service }} -p  {{ pypiserver.port }} {{ pypiserver.package_location }}
